# provisioning/templates/supervisor/uwsgi.conf.j2
[program:uwsgi]
command={{ virtualenv_path }}/bin/uwsgi --ini {{ uwsgi.ini_file}}
#environment=DJANGO_SETTINGS_MODULE={{ django_settings_module }}
