# provisioning/roles/supervisor/templates/supervisor/uwsgi.conf.j2
[program:uwsgi]
command={{ virtualenv_path }}/bin/uwsgi --ini {{ uwsgi_ini_file}}
priority=997
